# labs/realistic/vps/proxy-rs-client.service
[Unit]
Description=proxy-rs client
After=network-online.target
Wants=network-online.target

[Service]
User=proxy-rs
Group=proxy-rs
ExecStart=/usr/local/bin/proxy-rs run -c /etc/proxy-rs/client.json
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/proxy-rs

[Install]
WantedBy=multi-user.target
